# debian/postinst
#!/bin/sh -e

PKG="ds4drv"

if [ -e /etc/init.d/ds4drv ]; then
  update-rc.d -f ds4drv defaults
  update-rc.d ds4drv enable
fi

#DEBHELPER
