# debian/prerm
#!/bin/bash -e

# disable the service when we uninstall
systemctl enable ds4drv.service
